# daemon/script/hera-daemon.service
[Unit]
Description=Daemon of Hera
After=network.target

[Service]
Type=simple
PIDFile=/run/hera-daemon.pid
ExecStart=/usr/local/bin/hera-daemon /var/hera/data/ /var/hera/profiles.json /var/hera/logs/hera-daemon > /dev/null 2>&1 & 
Restart=on-failure
RestartSec=42s
User=root
Group=root

[Install]
WantedBy=multi-user.target
